# bftools/scripts/systemd/netidx@.service
[Unit]
Description=Netidx Services as %i

[Service]
User=%i
ExecStart=/usr/local/bin/netidx activation -f

[Install]
WantedBy=default.target
